# cc-heel-on-resume/cc-heel-on-resume.service
[Unit]
Description=Heel Claude Code /remote-control on every resume from suspend
Documentation=https://github.com/anthropics/claude-code/issues/34255
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/cc-heel-watcher
Restart=always
RestartSec=5
# dbus-monitor exits if the system bus restarts; the watcher restarts
# itself with a 5s backoff which is enough for transient bus blips.

# Hardening: read-only filesystem view except for the user runtime dir
# and a couple of obvious user paths the heel script touches transitively.
#
# PrivateTmp MUST be false: tmux's socket lives at /tmp/tmux-$UID/default
# (the real /tmp), NOT in %t/XDG_RUNTIME_DIR as originally assumed. With
# PrivateTmp=true the service got an isolated empty /tmp and every run
# failed at `tmux list-sessions` with "no tmux server, exit" — silently
# no-op'ing on every resume. ProtectSystem=strict keeps /tmp read-only as
# a mount, but connecting to the unix socket does not require a writable
# mount, so list + send-keys both work. (Verified 2026-06-11.)
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t %h/.local/bin
PrivateTmp=false
NoNewPrivileges=true

[Install]
WantedBy=default.target
